# freshports-freshports/files/freshports.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: freshports
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable freshports:
# freshports_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable freshports
#

. /etc/rc.subr

name="freshports"
rcvar=${name}_enable

pidfile="/var/run/${name}/${name}.pid"
freshports_user="freshports"
freshports_command="%%FP_SERVICESDIR%%/freshports.sh"
command="/usr/sbin/daemon"

load_rc_config $name

: ${freshports_enable:=NO}
: ${freshports_user:=freshports}
: ${freshports_group:=freshports}
: ${freshports_syslog_facility:=local5}

required_files="%%ETCDIR%%/freshports.sh"

start_precmd=freshports_prestart

freshports_prestart()
{
	# create the file pid, and directory, with correct permissions
	if [ ! -e ${pidfile} ]; then
		install -o ${freshports_user} -g ${freshports_group} /dev/null ${pidfile};
	else
		chown ${freshports_user}:${freshports_group} ${pidfile};
	fi

	# create the signals directory
	if [ ! -d "/var/db/${freshports_user}/signals" ]; then
		echo creating "/var/db/${freshports_user}/signals"
		install -d -o ${freshports_user} -g ${freshports_group} "/var/db/${freshports_user}/signals"
	fi

	# create the message-queues directory
	if [ ! -d /var/db/${freshports_user}/message-queues ]; then
		echo creating /var/db/${freshports_user}/message-queues
		install -d -o ${freshports_user} -g ${freshports_group} -m 0775 /var/db/${freshports_user}/message-queues
	fi

        # required dirs
        REQUIRED_DIRS="archive incoming recent retry"
        for dir in $REQUIRED_DIRS
        do
		# create the incoming directory
		if [ ! -d /var/db/${freshports_user}/message-queues/$dir ]; then
			echo creating /var/db/${freshports_user}/message-queues/$dir
			install -d -o ${freshports_user} -g ${freshports_group} -m 0775 /var/db/${freshports_user}/message-queues/$dir
		fi
	done
}

command_args="-P ${pidfile} -t ${name} -T ${name} -l ${freshports_syslog_facility} ${freshports_command}"

run_rc_command "$1"
