# freshports-freshports/pkg-plist
@sample(root,freshports,640) %%ETCDIR%%/freshports.sh.sample

%%FP_SERVICESDIR%%/freshports.sh

etc/rc.d/freshports
@sample etc/syslog.d/freshports-daemon.conf.sample

@dir(freshports,freshports,) /var/run/freshports
@dir(root,wheel,755) %%FP_SERVICESDIR%%
@dir(freshports,freshports,775) %%LOGDIR%%
