# freshports-ingress/pkg-plist
@sample(root,ingress,640) %%ETCDIR%%/ingress.sh.sample

%%FP_SERVICESDIR%%/ingress.sh

etc/rc.d/ingress
@sample etc/syslog.d/ingress-daemon.conf.sample
@sample etc/newsyslog.conf.d/ingress.sample

@dir(ingress,ingress,) /var/run/ingress
@dir(root,wheel,755) %%FP_SERVICESDIR%%
@dir(freshports,freshports,775) %%LOGDIR%%
