# py-freshports-fp-listen/pkg-plist
%%ETCDIR%%/fp-listen.ini.sample
etc/rc.d/fp_listen
%%PYTHON_SITELIBDIR%%/fp-listen/__pycache__/fp-listen%%PYTHON_TAG%%.opt-1.pyc
%%PYTHON_SITELIBDIR%%/fp-listen/__pycache__/fp-listen%%PYTHON_TAG%%.pyc
%%PYTHON_SITELIBDIR%%/fp-listen/fp-listen.py
%%PYTHON_SITELIBDIR%%/fp_listen-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/PKG-INFO
%%PYTHON_SITELIBDIR%%/fp_listen-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/SOURCES.txt
%%PYTHON_SITELIBDIR%%/fp_listen-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/dependency_links.txt
%%PYTHON_SITELIBDIR%%/fp_listen-%%PORTVERSION%%-py%%PYTHON_VER%%.egg-info/top_level.txt
